# provision/ansible/roles/podman_jellyfin/templates/jellyfin.timer.j2
{{ ansible_managed|comment }}

[Timer]
OnBootSec=20sec
