# provision/ansible/roles/unique-container/files/etc/systemd/mnt-usb.automount
[Unit]
Description=USB Automount

[Automount]
Where=/mnt/usb

[Install]
WantedBy=multi-user.target
