# provision/ansible/roles/podman-pihole/templates/pihole.service.j2
[Unit]
Description=Pi-hole
Documentation=https://pi-hole.net/
Wants=network.target
After=network.target
StartLimitIntervalSec=10

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
RestartSec=10
StartLimitBurst=3
# docker.io rate-limit delays can abort startup if timeout is too short.
TimeoutStartSec=300
TimeoutStopSec=70
ExecStartPre=/bin/rm --force %t/container-pihole.pid %t/container-pihole.ctr-id
ExecStartPre=-/usr/bin/podman pull docker.io/pihole/pihole:{{ pihole_version }}
ExecStart=/usr/bin/podman run \
    --cgroups no-conmon \
    --cidfile %t/container-pihole.ctr-id \
    --conmon-pidfile %t/container-pihole.pid \
    --detach \
    #--dns=127.0.0.1 \
    #--dns=1.1.1.1 \
    --env "FTLCONF_LOCAL_IPV4=192.168.178.52" \
    #--env "INTERFACE=eno1" \
    --env "TZ={{ config_system_timezone }}" \
    --env "WEBPASSWORD={{ pihole_webpassword }}" \
    --hostname %n \
    --label "traefik.enable=true" \
    --label "traefik.http.routers.pihole.entrypoints=websecure" \
    --label "traefik.http.routers.pihole.middlewares=oauth@file" \
    --label "traefik.http.routers.pihole.rule=Host(`%n.{{ inventory_hostname }}.{{ realm }}`)" \
    --label "traefik.http.routers.pihole.tls.certresolver={{ realm }}" \
    --label "traefik.http.routers.pihole.tls.domains[0].main=%n.{{ inventory_hostname }}.{{ realm }}" \
    --label "traefik.http.services.pihole.loadbalancer.server.port=80" \
    --log-driver journald \
    --mount "type=bind,source=/etc/ssl/certs/ca-certificates.crt,target=/etc/ssl/certs/ca-certificates.crt,ro=true" \
    --mount "type=bind,source={{ pihole_basedir }}/etc/pihole,target=/etc/pihole" \
    --mount "type=bind,source={{ pihole_basedir }}/etc/dnsmasq.d,target=/etc/dnsmasq.d" \
    --name %n \
    --network "service" \
    --publish "192.168.178.52:53:53/tcp" \
    --publish "192.168.178.52:53:53/udp" \
    --publish "8053:80/tcp" \
    --pull never \
    --replace \
    --rm \
    pihole:{{ pihole_version }}

ExecStop=/usr/bin/podman stop --ignore --cidfile %t/container-pihole.ctr-id --time 10
ExecStopPost=/usr/bin/podman rm --ignore --force --cidfile %t/container-pihole.ctr-id
PIDFile=%t/container-pihole.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
