# provision/ansible/roles/docker-cloudflare-dns/templates/cloudflare-dns.timer.j2
[Unit]
Description=Logs some system statistics to the systemd journal
Requires=cloudflare-dns.service

[Timer]
Unit=cloudflare-dns.service
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
